# hello/robot-hello.service.example
[Unit]
Description=GR6 hello (example child service)
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/share/GR6-v2
ExecStart=/usr/bin/python3 hello/app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
